# crontab
# Run homebase_calendar_sync once daily at midnight
0 0 * * * /usr/local/bin/homebase_calendar_sync
